# services/fpm/run
#!/bin/bash

exec /usr/sbin/php-fpm7.0 -F --allow-to-run-as-root
